# detections/siem/splunk/savedsearches.conf
# savedsearches.conf — sample compiled detections for Splunk.
#
# These are the DEPLOYABLE form of a handful of rules in ../../sigma, hand-compiled
# here to show the "Sigma is source of truth -> compile to your backend" step from
# DEFENSE-METHODOLOGY.md. In a real pipeline you would generate these with
# `sigma convert -t splunk ../../sigma/...` and drop the output in; this file is the
# worked example. Field names follow the Splunk Windows TA schema — adjust to your
# CIM/normalization. Drop into $SPLUNK_HOME/etc/apps/<app>/local/.
#
# Each search mirrors a PURPLE-TEAM.md block in dotfiles-Kali; validate by running
# the paired attack (htpx) and confirming the alert fires.

[Kerberoasting - RC4 service ticket (T1558.003)]
search = index=main EventCode=4769 Service_Name!="*$" Service_Name!="krbtgt" Ticket_Encryption_Type=0x17 \
| stats dc(Service_Name) AS ServiceAccounts values(Service_Name) by Client_Address, Account_Name \
| sort -ServiceAccounts
description = RC4 TGS for a user SPN — the Kerberoast downgrade tell. Sigma: credential_access/kerberoasting_rc4_tgs.yml
dispatch.earliest_time = -15m
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4
action.notable = 1

[Password spray - one source, many accounts (T1110.003)]
search = index=main EventCode=4625 NOT (Source_Network_Address IN ("-","127.0.0.1","::1")) \
| eval Account=mvindex(Account_Name,1) \
| stats dc(Account) AS Accounts by host, Source_Network_Address \
| where Accounts > 10 | sort -Accounts
description = One source failing against many distinct accounts. Sigma: credential_access/password_spray_4625.yml
dispatch.earliest_time = -15m
cron_schedule = */10 * * * *
enableSched = 1
alert.severity = 4
action.notable = 1

[DCSync - replication right by a non-DC (T1003.006)]
search = index=main EventCode=4662 Access_Mask="0x100" Security_ID!="S-1-5-18" \
| stats count by host, Account_Name, Object_Server | sort -count
description = Directory-replication right requested by a non-system principal. Sigma: credential_access/dcsync_replication_4662.yml
dispatch.earliest_time = -60m
cron_schedule = 0 * * * *
enableSched = 1
alert.severity = 5
action.notable = 1

[Scheduled task created with suspicious action (T1053.005)]
search = index=main EventCode=4698 \
| regex Task_Content="(?i)(-enc\b|-w\s+hidden|FromBase64|\\\\Users\\\\|\\\\Temp\\\\|mshta|regsvr32|rundll32|powershell.*(http|iex))" \
| table _time, host, Subject_Account_Name, Task_Name, Task_Content
description = Persistence via schtasks with an encoded/hidden/LOLBin action. Sigma: persistence/scheduled_task_suspicious_4698.yml
dispatch.earliest_time = -60m
cron_schedule = 5 * * * *
enableSched = 1
alert.severity = 4
action.notable = 1

[WMI subscription persistence - command-line consumer (T1546.003)]
search = index=main (EventCode=20 OR EventCode=21) \
| regex Destination="(?i)(powershell|cmd\.exe|mshta|wscript|cscript|rundll32|-enc|FromBase64)" \
| table _time, host, User, Name, Type, Destination, Query
description = Fileless WMI event-consumer persistence (Sysmon 19/20/21). Sigma: persistence/wmi_event_subscription_consumer.yml
dispatch.earliest_time = -60m
cron_schedule = 10 * * * *
enableSched = 1
alert.severity = 5
action.notable = 1
