# detections/sigma/persistence/scheduled_task_suspicious_4698.yml
title: Scheduled task created with a suspicious action
id: 8ae44ad2-6331-4444-97f4-b2eac893717c
status: experimental
description: |
  Detects scheduled-task persistence: a 4698 task-created event whose task XML runs
  encoded/hidden PowerShell, a LOLBin (mshta/regsvr32/rundll32), or something from a
  user-writable/temp path. Detect on the action, not the task name; 4702 catches the
  modify-an-existing-task variant. Baseline your software's legit tasks and alert on
  the rest.
  Validate (purple): run schtasks /create (or nxc -M schtask_as) from dotfiles-Kali
  (hacktheplanet "Persistence" fold; htpx pair schtask-persist). Requires the
  Object Access > Other Object Access Events audit subcategory.
references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://github.com/Gerrrt/htpx/blob/main/entries/blue/schtask-4698.md
author: dotfiles-Defense
date: 2026-06-30
tags:
  - attack.persistence
  - attack.execution
  - attack.t1053.005
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4698
  suspicious:
    TaskContent|contains:
      - '-enc'
      - '-EncodedCommand'
      - '-w hidden'
      - '-WindowStyle Hidden'
      - 'FromBase64String'
      - '\Users\'
      - '\Temp\'
      - '\mshta.exe'
      - '\regsvr32.exe'
      - '\rundll32.exe'
      - 'DownloadString'
      - 'IEX'
  condition: selection and suspicious
falsepositives:
  - Software installers that register tasks invoking helper scripts (allowlist).
level: high
